# profile.d/25pager.sh
# Use less for paging.
export PAGER="less"
