# profile.d/60wine.sh
# Suppress all fixme level log messages.
export WINEDEBUG=fixme-all
